# trunk/epispider.cron
0 * * * * wget http://epispider.berlios.de/class.promed.php > /dev/null
